# rootfs/etc/s6-overlay/s6-rc.d/sbfspotupload/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the SBFspotUpload
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

CONFIG_PATH=/data/options.json


# ---- Symbolic Links ----
ln -sf /usr/bin/sbfspot/SBFspot /SBFspot
ln -sf /usr/bin/sbfspot/SBFspot.cfg /SBFspot.cfg
ln -sf /usr/bin/sbfspot/SBFspotUpload.cfg /SBFspotUpload.cfg

# SBFspotUpload
# Enable SBFspotUpload
var="$(jq --raw-output '.Upload' $CONFIG_PATH)"
if bashio::var.equals "${var}" "true"; then
   bashio::log.cyan '|   Starting SBFspotUpload  |'
   exec /usr/bin/sbfspot/SBFspotUploadDaemon -c /usr/bin/sbfspot/SBFspotUpload.cfg
else
   bashio::log.yellow '|    SBFspotUploadDaemon Not Activated    |'
   bashio::log.yellow '|     change your options selection       |'
   bashio::log.yellow '|             to enable it                |'
fi
